# airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root --login-program /usr/local/bin/autostart_installer.sh - $TERM
Type=idle
